# .config/systemd/user/mpd.socket
[Unit]
Description=Music Player Daemon

[Socket]
ListenStream=/var/run/user/1000/mpd
ListenStream=6600

[Install]
WantedBy=sockets.target
