# .crontab
*/3 * * * * /home/kyle/.bin/mailrun
* * * * * /home/kyle/.bin/newsbeuter-update
