# .crontab
*/3 * * * * /home/kyle/.local/bin/mailrun
*/3 * * * * /home/kyle/.local/bin/crontab-update

#vim: set ft=crontab ts=8 sw=4 tw=140 et :
