# .crontab
*/3 * * * * /home/kyle/.bin/mailrun
*/3 * * * * /home/kyle/.bin/newsbeuter-update
#*/3 * * * * /home/kyle/.bin/checkip
#*/3 * * * * /home/kyle/.bin/anime-sync
*/3 * * * * /home/kyle/.bin/crontab-update

#vim: set ft=crontab ts=8 sw=4 tw=140 et :
